# packaging/deb/debian/postinst
#!/bin/sh
# Prism post-installation script
# This script is run after the package is installed

set -e

# Source debconf library for configuration management
. /usr/share/debconf/confmodule

# Package name
PACKAGE="prism"
USER="prism"
GROUP="prism"
HOME_DIR="/var/lib/prism"
LOG_DIR="/var/log/prism"
CONFIG_DIR="/etc/prism"

case "$1" in
    configure)
        # Create prism system user and group if they don't exist
        if ! getent group "$GROUP" >/dev/null 2>&1; then
            addgroup --system "$GROUP"
        fi
        
        if ! getent passwd "$USER" >/dev/null 2>&1; then
            adduser --system --ingroup "$GROUP" --home "$HOME_DIR" \
                    --no-create-home --disabled-login \
                    --shell /usr/sbin/nologin \
                    --gecos "Prism System User" "$USER"
        fi
        
        # Create and set ownership of directories
        mkdir -p "$HOME_DIR" "$LOG_DIR" "$CONFIG_DIR/aws"
        
        # Set proper ownership
        chown "$USER:$GROUP" "$HOME_DIR"
        chown "$USER:$GROUP" "$LOG_DIR"
        chown root:"$GROUP" "$CONFIG_DIR"
        chown root:"$GROUP" "$CONFIG_DIR/aws"
        
        # Set proper permissions
        chmod 755 "$HOME_DIR"
        chmod 750 "$LOG_DIR"
        chmod 750 "$CONFIG_DIR"
        chmod 750 "$CONFIG_DIR/aws"
        
        # Set permissions on config files
        if [ -f "$CONFIG_DIR/daemon.conf" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/daemon.conf"
            chmod 640 "$CONFIG_DIR/daemon.conf"
        fi
        
        if [ -f "$CONFIG_DIR/aws/config.template" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/aws/config.template"
            chmod 640 "$CONFIG_DIR/aws/config.template"
        fi
        
        if [ -f "$CONFIG_DIR/aws/credentials.template" ]; then
            chown root:"$GROUP" "$CONFIG_DIR/aws/credentials.template"
            chmod 640 "$CONFIG_DIR/aws/credentials.template"
        fi
        
        # Enable systemd service
        if [ -d /run/systemd/system ]; then
            systemctl --system daemon-reload >/dev/null || true
            if ! systemctl is-enabled prism.service >/dev/null 2>&1; then
                systemctl enable prism.service >/dev/null || true
            fi
            
            # Create systemd drop-in directory for local overrides
            mkdir -p /etc/systemd/system/prism.service.d
        fi
        
        # Display installation success message
        cat << 'EOF'

╔══════════════════════════════════════════════════════════════════════════════╗
║                   Prism Successfully Installed                    ║
╠══════════════════════════════════════════════════════════════════════════════╣
║                                                                              ║
║  🎉 Prism v0.4.2 has been installed successfully!                ║
║                                                                              ║
║  📋 What's installed:                                                        ║
║     ✅ Core binaries: cws, prismd                                              ║
║     ✅ Systemd service: prism.service                             ║
║     ✅ Configuration: /etc/prism/                                 ║
║     ✅ User account: prism (system user)                          ║
║                                                                              ║
║  🚀 Next steps:                                                              ║
║                                                                              ║
║     1. Configure AWS credentials:                                            ║
║        sudo cp /etc/prism/aws/credentials.template \             ║
║               /etc/prism/aws/credentials                          ║
║        sudo cp /etc/prism/aws/config.template \                  ║
║               /etc/prism/aws/config                               ║
║        sudo nano /etc/prism/aws/credentials                      ║
║                                                                              ║
║     2. Start the service:                                                    ║
║        sudo systemctl start prism                                ║
║                                                                              ║
║     3. Test the installation:                                                ║
║        cws --version                                                         ║
║        cws templates                                                         ║
║                                                                              ║
║  📚 Documentation: https://github.com/scttfrdmn/prism            ║
║  🐛 Issues: https://github.com/scttfrdmn/prism/issues            ║
║                                                                              ║
║  💡 Quick start guide:                                                       ║
║     Run 'man cws' for detailed usage instructions                           ║
║     Run 'man prismd' for daemon configuration help                            ║
║                                                                              ║
╚══════════════════════════════════════════════════════════════════════════════╝

EOF
        ;;
        
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
        
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# Automatically added by dh_systemd_enable/11.1.6ubuntu2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
    # This will only remove masks created by d-s-h on package removal.
    deb-systemd-helper unmask 'prism.service' >/dev/null || true

    # was-enabled defaults to true, so new installations run enable.
    if deb-systemd-helper --quiet was-enabled 'prism.service'; then
        # Enables the unit on first installation, creates new
        # symlinks on upgrades if the unit file has changed.
        deb-systemd-helper enable 'prism.service' >/dev/null || true
    else
        # Update the statefile to add new symlinks (if any), which need to be
        # cleaned up on purge. Also remove old symlinks.
        deb-systemd-helper update-state 'prism.service' >/dev/null || true
    fi
fi
# End automatically added section

exit 0